# systemd/ptpping-capture.service
[Unit]
Description=PTPPing Audio Capture
After=network.target ptp4l.service influxdb.service
Wants=ptp4l.service influxdb.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /opt/ptpping/ptpping.py --role capture --config /etc/ptpping/config.toml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Audio device access
SupplementaryGroups=audio

# PTP time access
SupplementaryGroups=ptp

[Install]
WantedBy=multi-user.target
